# supervisor/publication-resend-redpack.ini
[program:resend-redpack]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work --queue=resend_redpack --tries=1
autostart=true
autorestart=true
numprocs=1
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/resend_redpack_queue.log