# supervisor/publication-demand.ini
[program:demand]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work --queue=demand --tries=3 --memory=128
autostart=true
autorestart=true
numprocs=5
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/queue_demand.log